# examples/gadget-dns-server.service
# Per-user (user) systemd: put this file in ~/.config/systemd/user/
# Then: systemctl --user daemon-reload
#       systemctl --user enable gadget-dns-server.service
#       systemctl --user start gadget-dns-server.service
# To run when not logged in: loginctl enable-linger $USER
#
# System-wide: put in /etc/systemd/system/ and use systemctl (no --user).

[Unit]
Description=Gadget DNS server (gadget endpoints, DoT/DoH/DoQ, ACME, DNSSEC)
Documentation=https://github.com/davidgroves/gadget-dns-server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/gadget-dns/bin/gadget-dns-server --config /home/gadget-dns/etc/config.yaml
Restart=on-failure
RestartSec=5

# If using privileged ports (e.g. 53, 443), run as root (default).
# For non-privileged ports, uncomment and set:
User=gadget-dns
Group=gadget-dns

# Optional: limit capabilities if not running as root but need port 53
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
# For per-user: use WantedBy=default.target instead of multi-user.target
WantedBy=multi-user.target
